# backup.service
[Unit]
Description=Restic backup from Ubuntu to Fedora

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup.sh
